# zk3/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper-server]
command=/etc/init.d/zookeeper-server start
