# hbase/hregionserver/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase-regionserver]
command=/etc/init.d/hbase-regionserver start
