# hmaster/supervisord.conf
[supervisord]
nodaemon=true

[program:hbase-master]
command=/etc/init.d/hbase-master start
