# hdfs/dn/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop-hdfs-datanode]
command=/etc/init.d/hadoop-hdfs-datanode start


[program:hadoop-yarn-nodemanager]
command=/etc/init.d/hadoop-yarn-nodemanager start
